# crontab/crontab.sh
#! /bin/sh
/usr/bin/php `dirname $0`/update.php