# crontab/crontab.sh
#! /bin/sh
/usr/bin/php /var/www/html/stock/crontab/update.php
